# flatpak-authenticator.service.in
[Unit]
Description=flatpak authenticator

[Service]
BusName=org.flatpak.Authenticator
ExecStart=@libexecdir@/flatpak-authenticator
Type=dbus
